# ssh-ubuntu/config/supervisord/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:ssh]
command = /usr/sbin/sshd -D -e
priority = 1
autostart = true
autorestart = true
